# build-agents/freebsd-build-agent-image/idle-check.crontab
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

* * * * * root /usr/local/sbin/idle-check
