# build-infrastructure/git-cache-server-image/update-git-mirrors.cron
*/30 * * * * git /usr/local/bin/update-git-mirrors
